# packages/net-proxy/tinyproxy/files/systemd/tinyproxy.service
[Unit]
Description = Startup script for the tinyproxy server
After = network.target

[Service]
Type = forking
PIDFile = /run/tinyproxy.pid
ExecStart = /usr/sbin/tinyproxy -c /etc/tinyproxy/tinyproxy.conf
ExecReload = /bin/kill -HUP $MAINPID
KillMode = process
